# templates/minecraft.service.j2
[Unit]
Description=Minecraft Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/minecraft
ExecStart=/opt/minecraft/start.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
